# infra/fly/macmini.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchDaemon for the Understudy Mac Mini browser pool.
  architecture.md §12: "Mac Mini pool runs TinyFish browsers; both [Fly + Mac
  Mini] verify SLSA L2 + cosign on boot."

  Install:
    sudo cp infra/fly/macmini.plist /Library/LaunchDaemons/com.understudy.tinyfish.plist
    sudo launchctl bootstrap system /Library/LaunchDaemons/com.understudy.tinyfish.plist

  Uninstall:
    sudo launchctl bootout system /Library/LaunchDaemons/com.understudy.tinyfish.plist

  macmini-start.sh MUST succeed before TinyFish CLI runs. KeepAlive=true means
  launchd respawns on crash, but a non-zero exit from the verify step loops
  through the KeepAlive backoff — intentional: we would rather hot-loop a
  verify check than silently boot an unverified image.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.understudy.tinyfish</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/macmini-start.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/var/log/understudy/tinyfish.out.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/understudy/tinyfish.err.log</string>

    <key>WorkingDirectory</key>
    <string>/opt/understudy</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>IMAGE_REF</key>
        <string>ghcr.io/nihalnihalani/understudy-agent-base:latest</string>
        <key>COSIGN_EXPERIMENTAL</key>
        <string>1</string>
        <key>COSIGN_CERT_IDENTITY</key>
        <string>https://github.com/nihalnihalani/understudy/.github/workflows/release.yml@refs/heads/main</string>
        <key>COSIGN_CERT_OIDC_ISSUER</key>
        <string>https://token.actions.githubusercontent.com</string>
    </dict>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>
</dict>
</plist>
